# .github/tmp/package/js2/luci-app-parentcontrol/root/etc/hotplug.d/iface/97-parentcontrol
#!/bin/sh
. /lib/functions.sh
CONFIG=parentcontrol

[ "$ACTION" = ifup ] && [ "$(`uci -q get $CONFIG.@basic[0].enabled`)" == 1 ] && /etc/init.d/$CONFIG start

